# crontab
0 9 * * * /usr/local/bin/php /app/bin/console.php >> /var/log/cron.log 2>&1
